# Library/LaunchAgents/it.unbit.uwsgi.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <key>Label</key>
        <string>it.unbit.uwsgi</string>
        <key>ProgramArguments</key>
        <array>
                <string>/Applications/MNPP/Library/uwsgi/bin/uwsgi</string>
                <string>--socket</string>
                <string>:9001</string>
                <string>--master</string>
                <string>--logto</string>
                <string>/Applications/MNPP/logs/uwsgi/uwsgi.log</string>
                <string>--vhost</string>
                <string>--no-site</string>
                <string>--pythonpath</string>
                <string>/Applications/MNPP/htdocs/</string>
                <string>--pythonpath</string>
                <string>/Applications/MNPP/Library/python26/lib/python2.6/lib-dynload/</string>
                <string>--pythonpath</string>
                <string>/Applications/MNPP/Library/python26/lib/python2.6/site-packages/</string>
                <string>--callable</string>
                <string>app</string>
                <string>--die-on-term</string>
        </array>
</dict>
</plist>